# wifi-auto.service
[Unit]
Description=WiFi Auto Connect
After=network.target
Wants=network.target

[Service]
Type=oneshot
ExecStart=/usr/sbin/wpa_supplicant -B -i wlan0 -c /etc/wpa_supplicant/wpa_supplicant.conf
ExecStartPost=/usr/sbin/dhcpcd -k wlan0 || true
ExecStartPost=/usr/sbin/dhcpcd -w wlan0 || true
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target